# dockerfiles/runit/proxmox-remove-subscription/finish
#!/bin/sh
# Simple finish script - just exit normally
exit 0