# dockerfiles/runit/proxmox-remove-subscription/run
#!/bin/bash
exec 2>&1

SERVICE_NAME="proxmox-remove-subscription"
MARKER_FILE="/var/lib/proxmox-backup/.subscription-nag-disabled"
SCRIPT_LOG="/var/log/proxmox-backup/subscription-removal.log"
PROXMOXLIB="/usr/share/javascript/proxmox-widget-toolkit/proxmoxlib.js"

echo "[$SERVICE_NAME] Starting subscription removal service..."

# Create log directory and file
mkdir -p /var/log/proxmox-backup
echo "=== Starting new run at $(date) ===" >> "$SCRIPT_LOG"

# ALWAYS create the correct script at runtime to avoid Docker cache issues
cat > /etc/proxmox-backup-default/pbs-post-install.sh << 'SCRIPT_END'
#!/bin/bash
# Don't use set -e because arithmetic operations can return non-zero

# Environment variables with defaults
PBS_ENTERPRISE=${PBS_ENTERPRISE:-"no"}
PBS_NO_SUBSCRIPTION=${PBS_NO_SUBSCRIPTION:-"yes"}
DISABLE_SUBSCRIPTION_NAG=${DISABLE_SUBSCRIPTION_NAG:-"yes"}

echo "PBS Post-Install Configuration:"
echo "  PBS_ENTERPRISE=${PBS_ENTERPRISE}"
echo "  PBS_NO_SUBSCRIPTION=${PBS_NO_SUBSCRIPTION}"
echo "  DISABLE_SUBSCRIPTION_NAG=${DISABLE_SUBSCRIPTION_NAG}"

VERSION="$(awk -F'=' '/^VERSION_CODENAME=/{ print $NF }' /etc/os-release)"

# Configure repositories
if [ "${PBS_ENTERPRISE}" = "no" ] || [ "${PBS_ENTERPRISE}" = "false" ]; then
    echo "Disabling enterprise repository..."
    echo "# deb https://enterprise.proxmox.com/debian/pbs ${VERSION} pbs-enterprise" > /etc/apt/sources.list.d/pbs-enterprise.list
fi

if [ "${PBS_NO_SUBSCRIPTION}" = "yes" ] || [ "${PBS_NO_SUBSCRIPTION}" = "true" ]; then
    echo "Enabling no-subscription repository..."
    echo "deb http://download.proxmox.com/debian/pbs ${VERSION} pbs-no-subscription" > /etc/apt/sources.list.d/pbs-install-repo.list
fi

# Disable subscription nag if requested
if [ "${DISABLE_SUBSCRIPTION_NAG}" = "yes" ] || [ "${DISABLE_SUBSCRIPTION_NAG}" = "true" ]; then
    echo "Disabling subscription nag in UI..."
    
    # Primary target file
    PROXMOXLIB="/usr/share/javascript/proxmox-widget-toolkit/proxmoxlib.js"
    
    if [ -f "$PROXMOXLIB" ]; then
        # Create original backup if not exists
        if [ ! -f "${PROXMOXLIB}.original" ]; then
            cp "$PROXMOXLIB" "${PROXMOXLIB}.original"
            echo "Created original backup: ${PROXMOXLIB}.original"
        fi
        
        echo "Applying subscription nag patches..."
        
        # Pattern 1: Fix the actual subscription check
        sed -i "s/res\.data\.status\.toLowerCase() !== 'active'/res.data.status.toLowerCase() === 'NoMoreNagging'/g" "$PROXMOXLIB"
        echo "✓ Patched subscription status check"
        
        # Pattern 2: Change the "No valid subscription" message
        sed -i "s/No valid subscription/Subscription OK/g" "$PROXMOXLIB"
        echo "✓ Patched subscription message"
        
        # Pattern 3: Alternative check pattern
        sed -i "s/data\.status !== 'active'/data.status === 'NoMoreNagging'/g" "$PROXMOXLIB"
        echo "✓ Patched alternative status check"
        
        # Pattern 4: Double-quote variant
        sed -i 's/data\.status\.toLowerCase() !== "active"/data.status.toLowerCase() === "NoMoreNagging"/g' "$PROXMOXLIB"
        echo "✓ Patched double-quote variant"
        
        # Verification
        echo ""
        echo "Verification Results:"
        if grep -q "NoMoreNagging" "$PROXMOXLIB"; then
            echo "✓ NoMoreNagging marker found - patch successful!"
        else
            echo "⚠ NoMoreNagging marker not found"
        fi
        
        if grep -q "No valid subscription" "$PROXMOXLIB"; then
            echo "⚠ Warning: 'No valid subscription' text still present"
        else
            echo "✓ 'No valid subscription' text successfully removed"
        fi
    else
        echo "ERROR: proxmoxlib.js not found"
        exit 1
    fi
    
    # Create marker file
    mkdir -p /var/lib/proxmox-backup
    touch /var/lib/proxmox-backup/.subscription-nag-disabled
fi

echo "Signaling PBS proxy to reload..."
pkill -HUP -f proxmox-backup-proxy 2>/dev/null || true

echo "PBS post-install configuration completed successfully!"
echo "NOTE: Clear your browser cache (Ctrl+F5) to see changes!"
exit 0
SCRIPT_END

chmod +x /etc/proxmox-backup-default/pbs-post-install.sh
echo "[$SERVICE_NAME] Script created/updated" | tee -a "$SCRIPT_LOG"

echo "[$SERVICE_NAME] Waiting for proxmoxlib.js to check current state..." | tee -a "$SCRIPT_LOG"

# Wait for proxmoxlib.js to exist first
for i in $(seq 1 60); do
    if [ -f "$PROXMOXLIB" ]; then
        echo "[$SERVICE_NAME] Found proxmoxlib.js at attempt $i" | tee -a "$SCRIPT_LOG"
        break
    fi
    sleep 2
done

# Check if patch is actually applied (not just marker)
if [ -f "$PROXMOXLIB" ]; then
    if grep -q "NoMoreNagging" "$PROXMOXLIB" 2>/dev/null; then
        echo "[$SERVICE_NAME] Patch already applied (NoMoreNagging found)" | tee -a "$SCRIPT_LOG"
        
        # Also check if the "No valid subscription" text is gone
        if grep -q "No valid subscription" "$PROXMOXLIB" 2>/dev/null; then
            echo "[$SERVICE_NAME] But 'No valid subscription' still exists - will reapply!" | tee -a "$SCRIPT_LOG"
            rm -f "$MARKER_FILE"
        else
            echo "[$SERVICE_NAME] All patches verified - no action needed" | tee -a "$SCRIPT_LOG"
            touch "$MARKER_FILE"
            exec sleep infinity
        fi
    else
        echo "[$SERVICE_NAME] Patch not applied - will apply now" | tee -a "$SCRIPT_LOG"
        rm -f "$MARKER_FILE"
    fi
fi

echo "[$SERVICE_NAME] Waiting for PBS services to start..." | tee -a "$SCRIPT_LOG"

# Wait for PBS binaries to exist
while [ ! -f "/usr/lib/x86_64-linux-gnu/proxmox-backup/proxmox-backup-api" ] && \
      [ ! -f "/usr/lib/aarch64-linux-gnu/proxmox-backup/proxmox-backup-api" ]; do
    echo "[$SERVICE_NAME] Waiting for PBS binaries..." | tee -a "$SCRIPT_LOG"
    sleep 2
done

# Wait for PBS API to be responsive
echo "[$SERVICE_NAME] Waiting for PBS API to respond..." | tee -a "$SCRIPT_LOG"
for i in $(seq 1 30); do
    if curl -k -s https://localhost:8007/ > /dev/null 2>&1; then
        echo "[$SERVICE_NAME] PBS API is ready at attempt $i" | tee -a "$SCRIPT_LOG"
        break
    fi
    sleep 2
done

# Give it extra time to ensure everything is fully loaded
echo "[$SERVICE_NAME] Waiting 10 seconds for services to stabilize..." | tee -a "$SCRIPT_LOG"
sleep 10

# Run the actual patching script
echo "[$SERVICE_NAME] Running subscription removal script..." | tee -a "$SCRIPT_LOG"
if /bin/bash /etc/proxmox-backup-default/pbs-post-install.sh >> "$SCRIPT_LOG" 2>&1; then
    echo "[$SERVICE_NAME] Script completed successfully" | tee -a "$SCRIPT_LOG"
    
    # Verify the patch was actually applied
    if grep -q "NoMoreNagging" "$PROXMOXLIB" 2>/dev/null; then
        echo "[$SERVICE_NAME] ✓ Patch verified - NoMoreNagging found" | tee -a "$SCRIPT_LOG"
        touch "$MARKER_FILE"
    else
        echo "[$SERVICE_NAME] ⚠ WARNING: Script ran but patch not verified!" | tee -a "$SCRIPT_LOG"
    fi
else
    EXIT_CODE=$?
    echo "[$SERVICE_NAME] Script failed with exit code $EXIT_CODE" | tee -a "$SCRIPT_LOG"
    echo "[$SERVICE_NAME] Last 20 lines of log:" | tee -a "$SCRIPT_LOG"
    tail -n 20 "$SCRIPT_LOG"
fi

echo "[$SERVICE_NAME] Service complete, keeping alive..." | tee -a "$SCRIPT_LOG"
exec sleep infinity