# dockerfiles/runit/proxmox-remove-subscription/run
#!/bin/sh
exec 2>&1

SERVICE_NAME="proxmox-remove-subscription"
MARKER_FILE="/var/lib/proxmox-backup/.subscription-nag-disabled"

echo "[$SERVICE_NAME] Starting service..."

# Check if already completed
if [ -f "$MARKER_FILE" ]; then
    echo "[$SERVICE_NAME] Subscription nag already removed, sleeping forever..."
    exec sleep infinity
fi

# Wait for PBS API to be ready (max 30 seconds)
echo "[$SERVICE_NAME] Waiting for PBS API to start..."
for i in $(seq 1 30); do
    if curl -k -s https://localhost:8007/ > /dev/null 2>&1; then
        echo "[$SERVICE_NAME] PBS API is ready"
        break
    fi
    sleep 1
done

# Run the post-install script
POST_INSTALL_SCRIPT="/etc/proxmox-backup-default/pbs-post-install.sh"
if [ -f "$POST_INSTALL_SCRIPT" ]; then
    echo "[$SERVICE_NAME] Running PBS post-install configuration..."
    if sh "$POST_INSTALL_SCRIPT"; then
        echo "[$SERVICE_NAME] Post-install script completed successfully"
        touch "$MARKER_FILE"
    else
        echo "[$SERVICE_NAME] WARNING: Post-install script failed with exit code $?"
    fi
else
    echo "[$SERVICE_NAME] Post-install script not found"
fi

# Sleep forever to prevent restarts
echo "[$SERVICE_NAME] Complete, sleeping forever..."
exec sleep infinity