# dockerfiles/runit/proxmox-remove-subscription/run
#!/bin/bash
exec 2>&1

SERVICE_NAME="proxmox-remove-subscription"
PROXMOXLIB="/usr/share/javascript/proxmox-widget-toolkit/proxmoxlib.js"

echo "[$SERVICE_NAME] Starting subscription removal service..."

# Wait for proxmoxlib.js to exist
for i in $(seq 1 60); do
    if [ -f "$PROXMOXLIB" ]; then
        echo "[$SERVICE_NAME] Found proxmoxlib.js"
        break
    fi
    sleep 2
done

# Apply patches if needed
if [ -f "$PROXMOXLIB" ]; then
    if ! grep -q "NoMoreNagging" "$PROXMOXLIB" 2>/dev/null; then
        echo "[$SERVICE_NAME] Applying patches..."
        
        # Backup original
        [ ! -f "${PROXMOXLIB}.original" ] && cp "$PROXMOXLIB" "${PROXMOXLIB}.original"
        
        # Apply all patches
        sed -i "s/res\.data\.status\.toLowerCase() !== 'active'/res.data.status.toLowerCase() === 'NoMoreNagging'/g" "$PROXMOXLIB"
        sed -i "s/No valid subscription/Subscription OK/g" "$PROXMOXLIB"
        sed -i "s/data\.status !== 'active'/data.status === 'NoMoreNagging'/g" "$PROXMOXLIB"
        sed -i 's/data\.status\.toLowerCase() !== "active"/data.status.toLowerCase() === "NoMoreNagging"/g' "$PROXMOXLIB"
        
        echo "[$SERVICE_NAME] Patches applied"
    else
        echo "[$SERVICE_NAME] Already patched"
    fi
fi

# DO NOT send any signals to the proxy - just sleep
echo "[$SERVICE_NAME] Complete"
exec sleep infinity